# rc.local
#!/usr/bin/env bash

iptables -A OUTPUT -p tcp -m tcp --dport 6697 -j REJECT --reject-with icmp-port-unreachable

